# debpkg/bin/etc/rc.local
#!/bin/bash
bash /etc/pvelive/network.sh &
bash /etc/pvelive/password.sh &
bash /etc/pvelive/hosts.sh &
exit 0
